# docker/supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/home/user/.supervisord.log
pidfile=/home/user/.supervisord.pid

# ---------------- Cursor (display :1, VNC 5901, noVNC 6080) ----------------
[program:xvnc-cursor]
command=/usr/local/bin/start-vnc.sh 1 5901
autorestart=true
priority=10
stdout_logfile=/home/user/.log-xvnc-cursor.log
stderr_logfile=/home/user/.log-xvnc-cursor.log

[program:cursor]
command=/usr/local/bin/start-app.sh 1 /opt/cursor/AppRun --no-sandbox --no-zygote --disable-gpu
autorestart=true
priority=20
stdout_logfile=/home/user/.log-cursor.log
stderr_logfile=/home/user/.log-cursor.log

[program:novnc-cursor]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5901
autorestart=true
priority=30
stdout_logfile=/home/user/.log-novnc-cursor.log
stderr_logfile=/home/user/.log-novnc-cursor.log

# ------------- Antigravity (display :2, VNC 5902, noVNC 6081) -------------
[program:xvnc-antigravity]
command=/usr/local/bin/start-vnc.sh 2 5902
autorestart=true
priority=10
stdout_logfile=/home/user/.log-xvnc-antigravity.log
stderr_logfile=/home/user/.log-xvnc-antigravity.log

# Use the real Electron binary directly, not /usr/bin/antigravity — that wrapper
# detects Docker-Desktop's WSL2 backend and blocks on an interactive [y/N] prompt.
[program:antigravity]
command=/usr/local/bin/start-app.sh 2 /usr/share/antigravity/antigravity --no-sandbox --disable-gpu
autorestart=true
priority=20
stdout_logfile=/home/user/.log-antigravity.log
stderr_logfile=/home/user/.log-antigravity.log

[program:novnc-antigravity]
command=/usr/bin/websockify --web=/usr/share/novnc 6081 localhost:5902
autorestart=true
priority=30
stdout_logfile=/home/user/.log-novnc-antigravity.log
stderr_logfile=/home/user/.log-novnc-antigravity.log

# ---------------- VS Code (code-server, localhost only) ----------------
[program:codeserver]
command=/usr/bin/code-server --auth none --bind-addr 127.0.0.1:8443 /home/user/project
autorestart=true
priority=30
environment=HOME="/home/user",USER="user"
stdout_logfile=/home/user/.log-codeserver.log
stderr_logfile=/home/user/.log-codeserver.log

# ---------------- Caddy: the single public entry on :8080 ----------------
[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autorestart=true
priority=40
environment=HOME="/home/user",XDG_CONFIG_HOME="/home/user/.config",XDG_DATA_HOME="/home/user/.local/share"
stdout_logfile=/home/user/.log-caddy.log
stderr_logfile=/home/user/.log-caddy.log
